# pokebot.service
[Unit]
Description=Str8Dex Pokemon Price Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/discord-pokemon-bot
ExecStart=/home/ubuntu/discord-pokemon-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
